# hanlon.sv.conf
[program:hanlon]
command=/bin/bash -c "/docker-entrypoint.sh"
redirect_stderr=true